# amr2mqtt/rootfs/etc/services.d/rtl_tcp/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Home Assistant Add-on: AMR2MQTT
# Take down the S6 supervision tree when AMR2MQTT fails
# ==============================================================================
if { s6-test ${1} -ne 0 }
if { s6-test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/service